# services/forward/run
#!/bin/sh
exec socat TCP-LISTEN:${DEBUG_PORT},fork,reuseaddr TCP:127.0.0.1:9223