# services/x11vnc/run
#!/bin/sh
exec /usr/bin/x11vnc -display :0 -nopw -forever -shared